# SystemdFiles/kindle-booklog.timer
[Unit]
Description=kindle-booklog

[Timer]
Unit=kindle-booklog.service
OnCalendar=*-*-* *:00,30:00

[Install]
WantedBy=timers.target